# stat.cron.sample
#!/bin/sh

set -e
export STAT_SERVICES=   # e.g. 'memcached mongodb nap nginx postfix'
/opt/src/stat/stat.bash /opt/src/mailer/mailer.bash

# --
